# services/hs-agent-main/run
#!/bin/sh
echo starting
export RAILS_ENV=production
exec 2>&1
exec /usr/lib/hs/agent/bin/HSAgent --main ................................................................................................................................................................
